# DEBIAN/postrm
#!/bin/bash

update-rc.d pjc-compmgr remove
