# DEBIAN/postinst
#!/bin/bash

update-rc.d pjc-compmgr defaults